# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[program:dbus]
command=/bin/bash -c "mkdir -p /var/run/dbus && dbus-daemon --system --nofork"
autorestart=true
priority=5

[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 1280x720x24
autorestart=true
priority=10

[program:xrdp-sesman]
command=/usr/sbin/xrdp-sesman --nodaemon
autorestart=true
priority=20

[program:xrdp]
command=/usr/sbin/xrdp --nodaemon
autorestart=true
priority=30

[program:websockify]
command=/bin/bash -c "/usr/bin/websockify --web=/usr/share/novnc 0.0.0.0:${PORT} 127.0.0.1:3390"
autorestart=true
priority=40
