# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[program:dbus]
command=/bin/bash -c "mkdir -p /var/run/dbus && dbus-daemon --system --nofork"
autorestart=true
priority=5

[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 1280x720x24
autorestart=true
priority=10
environment=DISPLAY=":1"

[program:xfce]
command=/bin/bash -c "export DISPLAY=:1 && startxfce4"
user=appuser
autorestart=true
priority=15
environment=DISPLAY=":1",HOME="/home/appuser",USER="appuser"

[program:x11vnc]
command=/usr/bin/x11vnc -display :1 -forever -shared -nopw -listen 127.0.0.1 -xkb
autorestart=true
priority=20
environment=DISPLAY=":1"

[program:websockify]
command=/bin/bash -c "/usr/bin/websockify --web=/usr/share/novnc 0.0.0.0:${PORT} 127.0.0.1:5900"
autorestart=true
priority=30
