# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:dbus]
command=/bin/bash -c "mkdir -p /var/run/dbus && dbus-daemon --system --nofork"
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 1280x720x24
autorestart=true
priority=10
environment=DISPLAY=":1"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:xfce]
command=/bin/bash -c "export DISPLAY=:1 && export HOME=/home/appuser && startxfce4"
user=appuser
autorestart=true
priority=15
environment=DISPLAY=":1",HOME="/home/appuser",USER="appuser"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:x11vnc]
command=/usr/bin/x11vnc -display :1 -forever -shared -nopw -listen 127.0.0.1 -xkb
autorestart=true
priority=20
environment=DISPLAY=":1"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:websockify]
command=/bin/bash -c "exec /usr/bin/websockify --web=/usr/share/novnc 0.0.0.0:${PORT} 127.0.0.1:5900"
autorestart=true
priority=30
environment=PORT="10000"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
